# .config/systemd/user/bw-serve.service
[Unit]
Description=Bitwarden CLI local REST API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%t/bw-session.env
ExecSearchPath=/usr/bin:/usr/sbin:/usr/local/bin:%h/.local/bin:%h/.nix-profile/bin
ExecStart=bw serve --port 8087 --hostname 127.0.0.1
Restart=on-failure
RestartSec=10

# Hardening
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=%h/.config/Bitwarden CLI

[Install]
WantedBy=default.target
